# deploy/systemd/g7mediabooster-inventory.service
[Unit]
Description=G7MediaBooster bounded non-destructive provider inventory audit
After=network-online.target g7mediabooster-api.service
Wants=network-online.target

[Service]
Type=oneshot
User=g7mediabooster
Group=g7mediabooster
ExecStart=/usr/local/bin/g7mb-worker --config /etc/g7mediabooster/g7mb.toml inventory
TimeoutStartSec=15min
UMask=0077

StateDirectory=g7mediabooster
RuntimeDirectory=g7mediabooster
ReadWritePaths=/var/lib/g7mediabooster /run/g7mediabooster
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

CPUQuota=25%
MemoryMax=192M
MemorySwapMax=0
TasksMax=24
